# make-deb/server/postinst
#!/bin/sh
set -e
case "$1" in
configure)
    if ! getent passwd tabularium >/dev/null; then
        adduser --system --group --disabled-login --shell /usr/sbin/nologin \
            --home /var/tabularium --no-create-home tabularium
    fi
    install -d -m700 -o tabularium -g tabularium /var/tabularium/data /var/tabularium/spill
    if [ ! -f /etc/tabularium/config.toml ]; then
        cp /etc/tabularium/config.toml.default /etc/tabularium/config.toml
    fi
    # Lock down operator config dir (Ferrum / archivion pattern).
    chown tabularium:tabularium /etc/tabularium
    chmod 700 /etc/tabularium
    chown tabularium:tabularium /etc/tabularium/config.toml.default
    chmod 600 /etc/tabularium/config.toml.default
    if [ -f /etc/tabularium/config.toml ]; then
        chown tabularium:tabularium /etc/tabularium/config.toml
        chmod 600 /etc/tabularium/config.toml
    fi
    if [ -x /bin/systemctl ] && [ -d /run/systemd/system ]; then
        systemctl daemon-reload || true
        systemctl enable tabularium-server.service || true
        if [ -f /run/tabularium-deb-was-active ]; then
            rm -f /run/tabularium-deb-was-active
            systemctl restart tabularium-server.service || true
        elif [ -n "$2" ]; then
            # Upgrade or reconfigure: restart only if already active (no surprise start).
            systemctl try-restart tabularium-server.service || true
        else
            systemctl start tabularium-server.service || true
        fi
    fi
    ;;
esac
exit 0
